# startup/fcserver.conf
description "Fadecandy controller"
author "Paul Showalter"
start on runlevel [2345]
script
	exec /usr/local/bin/fcserver /etc/fcconfig
end script
stop on shutdown
respawn
respawn limit 10 5
